# Web2py/contaminacion/cron/crontab
#crontab
11 0-23/1 * * * root python *applications/contaminacion/cron/calidadaire.py